# roles/ssl/templates/crontab
0 */12 * * * root AWS_CONFIG_FILE=/srv/certbot.aws.config /srv/certbot-auto renew --post-hook "service nginx restart"
